# systemd/scoreboard.service
[Unit]
Description=Raspberry Pi 5 Piomatter Baseball Scoreboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/scoreboard_3panel_pi5
ExecStart=/opt/scoreboard_3panel_pi5/scoreboard --backend piomatter --no-sudo -- --state-file /opt/scoreboard_3panel_pi5/scoreboard_state.json
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
